# files/supervisor.conf

{%- for service in pillar.supervisor.server.services %}
{%- if service.type == 'flower' %}

{%- if not "message_queue" in pillar.flower.server %}
{%- set service_config = pillar.flower.server %}
{%- else %}
{%- set service_config = pillar.flower.server.message_queue[service.name] %}
{%- endif %}

[program:{{ service.type }}_{{ service.name }}]
directory=/srv/flower
environment=PATH="/srv/flower/bin"
{%- if service_config.broker.engine == 'rabbitmq' %}

command=/srv/flower/bin/flower --port={{ service_config.bind.port }} --address={{ service_config.bind.address }} --broker=amqp://{{ service_config.broker.user }}:{{ service_config.broker.password }}@{{ service_config.broker.host }}:{{ service_config.broker.port }}/{{ service_config.broker.virtual_host }}{%- if service_config.identity is defined %} --keystone_auth={{ service_config.identity.host }}{%- endif %}

{%- elif service_config.broker.engine == 'redis' %}
command=/srv/flower/bin/flower --port={{ service_config.bind.port }} --address={{ service_config.bind.address }} --broker=redis://{{ service_config.broker.host }}:{{ service_config.broker.port }}/{{ service_config.broker.number }}{%- if service_config.identity is defined %} --keystone_auth={{ service_config.identity.host }}{%- endif %}

{%- endif %}
stdout_logfile=/srv/flower/logs/flower_access.log
stderr_logfile=/srv/flower/logs/flower_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

{%- endif %}
{%- endfor %}

{%- for service_name, service in pillar.supervisor.server.service.iteritems() %}
{%- if service.type == 'flower' %}

{%- if not "message_queue" in pillar.flower.server %}
{%- set service_config = pillar.flower.server %}
{%- else %}
{%- set service_config = pillar.flower.server.message_queue[service.name] %}
{%- endif %}

[program:{{ service.name }}]
directory=/srv/flower
environment=PATH="/srv/flower/bin"
{%- if service_config.broker.engine == 'rabbitmq' %}

command=/srv/flower/bin/flower --port={{ service_config.bind.port }} --address={{ service_config.bind.address }} --broker=amqp://{{ service_config.broker.user }}:{{ service_config.broker.password }}@{{ service_config.broker.host }}:{{ service_config.broker.port }}/{{ service_config.broker.virtual_host }}{%- if service_config.identity is defined %} --keystone_auth={{ service_config.identity.host }}{%- endif %}

{%- elif service_config.broker.engine == 'redis' %}
command=/srv/flower/bin/flower --port={{ service_config.bind.port }} --address={{ service_config.bind.address }} --broker=redis://{{ service_config.broker.host }}:{{ service_config.broker.port }}/{{ service_config.broker.number }}{%- if service_config.identity is defined %} --keystone_auth={{ service_config.identity.host }}{%- endif %}

{%- endif %}
stdout_logfile=/srv/flower/logs/flower_access.log
stderr_logfile=/srv/flower/logs/flower_error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600

{%- endif %}
{%- endfor %}
